# public/assets/common/files/etc/hotplug.d/iface/88-replace-apk-source
#!/bin/sh
[ "$ACTION" = "ifup" ] && [ "$INTERFACE" = "wan" ] && {
    sleep 30
    /usr/bin/replace-apk-source.sh
}
